# harvest/template/server/supervisor/development.ini
[program:{{ project_name }}-development]
directory = /home/devuser/webapps/{{ project_name }}-env/{{ project_name }}
command = /home/devuser/webapps/{{ project_name }}-env/bin/uwsgi --ini /home/devuser/webapps/{{ project_name }}-env/{{ project_name }}/server/uwsgi/development.ini
autostart = true
autorestart = true
